# cosmo-kinoite-cleanup-1.4/payload/etc/systemd/system/system-clean.timer
[Unit]
Description=Run extended system cleanup weekly

[Timer]
OnCalendar=Weekly
Persistent=true

[Install]
WantedBy=timers.target
